# databag.init
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org
 
START=10
STOP=15
 
start() {        
        echo start
        mkdir -p /tmp/databag/assets
        cp -n /opt/databag/databag.db /tmp/databag/
        /usr/bin/databag 7000 /tmp/databag/ > /tmp/databag/service.log
}                 
 
stop() {          
        echo stop
        killall databag 
}
