# databag.init
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org
 
START=10
STOP=15
 
start() {        
        echo start
        mkdir -p /var/lib/databag/assets
        cp -n /opt/databag/databag.db /var/lib/databag/
        /usr/bin/databag 7000 /opt/databag/web
}                 
 
stop() {          
        echo stop
        killall databag 
}
